# systemd/system/omarchy-firstboot.service
[Unit]
Description=Omarchy bootc first-boot root setup
Documentation=https://github.com/joshyorko/omarchy-bootc
After=local-fs.target
Before=greetd.service
# One-shot gate: script also writes this stamp to stay idempotent.
ConditionPathExists=!/var/lib/omarchy/.firstboot-done

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/omarchy/omarchy-setup.sh
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
